# app.service
[Unit]
Description=app

[Service]
ExecStart=/github/My_Project/app
Restart=always

[Install]
WantedBy=multi-user.target

